# curator4-995309161.conf
[program:curator4-995309161]
command=python3 /home/curatorN/curator4-995309161/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curatorN/curator4-995309161/test.err.log
stdout_logfile=/home/curatorN/curator4-995309161/test.out.log